# debian/statsd.upstart
description "statsd"
author      "rdio"

start on startup
stop on shutdown

env HOME=/
env CONFIG=/etc/statsd/config.js
env LOGFILE=/tmp/statsd.log
env USER=nobody

script

  if [ -e "/etc/default/statsd" ] ; then
    . /etc/default/statsd
  fi

  exec sudo -u $USER /usr/bin/node /usr/share/statsd/stats.js $CONFIG 2>&1 >>$LOGFILE

end script

